# dotfiles/launchagents/com.user.caenv.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.caenv</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>
      CERT="$HOME/personal/ansible-macos/certs/Cloudflare_CA.pem";
      launchctl setenv NODE_EXTRA_CA_CERTS "$CERT";
      launchctl setenv REQUESTS_CA_BUNDLE "$CERT";
      launchctl setenv SSL_CERT_FILE "$CERT";
      launchctl setenv CURL_CA_BUNDLE "$CERT"
    </string>
  </array>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
